# bin/wyze-mqtt.service
[2;2R]
Description=Wyzesense via Mqtt
Documentation=https://www.github.com/ozczecho
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/wyze-mqtt/bin
ExecStart=/opt/wyze-mqtt/bin/service.sh
Restart=always

[Install]
WantedBy=multi-user.target
